# systemd/ovpn-data-init.service
[Unit]
Description=Initialize OpenVPN panel data directory
Before=openvpn-ovpn.service stunnel-ovpn.service ovpn-panel.service
DefaultDependencies=no

[Service]
Type=oneshot
RemainAfterExit=yes
EnvironmentFile=-/etc/default/ovpn-panel
ExecStart=/opt/panel/scripts/ovpn-data-init.sh

[Install]
WantedBy=multi-user.target
